# valuescan-token-refresher.service
[Unit]
Description=ValueScan Token Auto Refresher
Documentation=https://github.com/ymylive/valuescan
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/valuescan/signal_monitor
ExecStart=/usr/bin/python3.9 /root/valuescan/signal_monitor/token_refresher.py --interval 0.8
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# 环境变量
Environment=DISPLAY=:99
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/root/valuescan/config/valuescan.env

[Install]
WantedBy=multi-user.target
