# valuescan-trader.service
[Unit]
Description=ValueScan Binance Trader (IPC)
Documentation=https://github.com/ymylive/valuescan
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/valuescan/binance_trader
ExecStart=/usr/bin/python3.9 /root/valuescan/binance_trader/ipc_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=valuescan-trader

Environment=PYTHONUNBUFFERED=1
Environment=VALUESCAN_TRADER_MODE=ipc

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

